# files/etc/cron.daily/letsencrypt-renew
#!/bin/bash
/usr/local/bin/letsencrypt renew --quiet --no-self-upgrade